# ionic-default/files/etc/hotplug.d/iface/99-ddnsupdate
#!/bin/sh
logger -t hotplug "Action: $ACTION Interface: $INTERFACE"
#if [ $ACTION == "ifup" ]; then
#	logger -t HOTPLUG "Device $INTERFACE started - updating DDNS"
#	/usr/bin/ddnsup
#fi
