# themes/luci-theme-gs/root/etc/uci-defaults/luci-theme-gs
#!/bin/sh
uci batch <<-EOF
	set luci.themes.Gs=/luci-static/gs
	set luci.main.mediaurlbase=/luci-static/gs
	commit luci
EOF
exit 0
