# themes/ionic/root/etc/uci-defaults/luci-theme-ionic
#!/bin/sh
uci batch <<-EOF
	set luci.themes.Ionic=/luci-static/ionic
	commit luci
EOF
exit 0
