# ionic-default/files/etc/uci-defaults/70-hostname
#!/bin/sh
MACADDR=`ifconfig eth0|grep HWaddr|sed -e 's/^.*HWaddr //g' -e 's/://g'`
uci batch <<-EOF
	set system.@system[0].hostname=ionic-$MACADDR
	commit system
EOF
exit 0
