# themes/thingmesh/root/etc/uci-defaults/luci-theme-thingmesh
#!/bin/sh
uci batch <<-EOF
	set luci.themes.ThingMesh=/luci-static/thingmesh
	commit luci
EOF
exit 0
